# onlitec-email-webhook.service
[Unit]
Description=Onlitec Email Security Webhook Listener
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/alfreire/onlitec-email-security
Environment="WEBHOOK_SECRET=onlitec-email-webhook-secret-2024"
ExecStart=/usr/bin/python3 /home/alfreire/onlitec-email-security/webhook-listener.py
Restart=always
RestartSec=10
StandardOutput=append:/home/alfreire/onlitec-email-security/webhook.log
StandardError=append:/home/alfreire/onlitec-email-security/webhook.log

[Install]
WantedBy=multi-user.target
